# debian/perdition-pbs.postrm
#! /bin/sh
# postrm script for perdition-pbs

#DEBHELPER#

if [ "$1" = "remove" ]; then
	update-rc.d -f perdition-pbs remove >/dev/null
fi

exit 0


